# retired_spells/moonbase/kde3/kdebase3/profile.d/kde3.rc
export  KDEDIR=/opt/lunar/kde/3
export  PATH=/opt/lunar/kde/3/bin:$PATH
